# crontab
0 5 * * * /home/username/bin/build_video_web >> /var/log/build_video_web.log 2>&1
